# setup/systemd_files/fan.service
[Unit]
Description=52Pi PWM Fan Control Service
After=multi-user.target

[Service]
Type=simple
# Use the system python3 interpreter as RPi.GPIO is installed via apt
ExecStart=/usr/bin/python3 /opt/ballbot/fan_control.py
Restart=always
# Optional: Add user/group if needed, though usually GPIO access requires group membership
# User=pi
# Group=gpio

[Install]
WantedBy=multi-user.target
