# third_party/dofbot-pro/New_Code/Orin/DofbotPro_source/dofbot_pro/Yahboom_Oled/yahboom_oled.service
[Unit]
Description=yahboom_oled
After=multi-user.target

[Service]
Type=idle
User=jetson
ExecStart=/bin/sh -c "python3 /home/jetson/dofbot_pro/Yahboom_Oled/yahboom_oled.py"
WorkingDirectory=/home/jetson/

[Install]
WantedBy=multi-user.target
