# third_party/dofbot-pro/New_Code/Orin/DofbotPro_source/software/oled_yahboom/yahboom_oled.service
[Unit]
Description=yahboom_oled start service
After=multi-user.target

[Service]
Type=idle
User=jetson
ExecStart=/bin/sh -c "python3 /home/jetson/software/oled_yahboom/yahboom_oled.py"
WorkingDirectory=/home/jetson

[Install]
WantedBy=multi-user.target
